# roles/raspi_base/templates/lte-policy-routing.service.j2
[Unit]
Description=Apply LTE source-based policy routing for static public IP
After=NetworkManager.service ModemManager.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/fix-lte-policy-routing.sh
RemainAfterExit=yes
TimeoutStartSec=240

[Install]
WantedBy=multi-user.target
